# crontab
*/4 * * * * /usr/bin/python3 /app/runner.py /app/config
