# crontab
*/4 * * * * /usr/bin/python /app/runner.py /app/config
